# root/etc/s6-overlay/s6-rc.d/init-backup-on-start/up
with-contenv bash /etc/s6-overlay/s6-rc.d/init-backup-on-start/run.sh