# root/etc/s6-overlay/s6-rc.d/http-api/run
#!/command/execlineb -P
/bin/bash /app/httpd.sh start